# user/default/files/etc/hotplug.d/iface/21-domain
#!/bin/bash
echo "address=/op.cn/`uci show network.lan.ipaddr | awk -F"'" '{print $2}'`" > /tmp/dnsmasq.d/domain.conf
/etc/init.d/dnsmasq restart
